# deploy/chatbot-python.conf
; ============================================================
; Supervisor config untuk Python Flask ChatbotRSJ
; Upload file ini ke server: /etc/supervisor/conf.d/chatbot-python.conf
; Lalu jalankan: sudo supervisorctl reread && sudo supervisorctl update
; ============================================================

[program:chatbot-python]
command=/var/www/chatbot/.venv/bin/python -m python_chatbot.bot_ai.app
directory=/var/www/chatbot
user=www-data
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stderr_logfile=/var/log/supervisor/chatbot-python.err.log
stdout_logfile=/var/log/supervisor/chatbot-python.out.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
environment=
    FLASK_ENV="production",
    PYTHONPATH="/var/www/chatbot",
    HF_HUB_OFFLINE="1",
    TRANSFORMERS_OFFLINE="1"
